# deploy/launchd/com.itpe.splitter.mlx.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.itpe.splitter.mlx</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/turtlesoup0-macmini/Projects/itpe-topic-splitter/venv/bin/python3</string>
        <string>-m</string>
        <string>mlx_lm</string>
        <string>server</string>
        <string>--model</string>
        <string>Jiunsong/supergemma4-26b-uncensored-mlx-4bit-v2</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8090</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/turtlesoup0-macmini/Projects/itpe-topic-splitter</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/turtlesoup0-macmini/Library/Logs/itpe-splitter/mlx.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/turtlesoup0-macmini/Library/Logs/itpe-splitter/mlx.err.log</string>
</dict>
</plist>
