# deploy/launchd/com.itpe.splitter.cloudflared.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.itpe.splitter.cloudflared</string>

    <!-- Named Tunnel: topic-splitter.tech-insight.org -->
    <!-- 설정: ~/.cloudflared/config.yml -->
    <!-- 기동마다 URL 변경되는 trycloudflare 임시 터널에서 고정 도메인 전환 -->
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/cloudflared</string>
        <string>tunnel</string>
        <string>run</string>
        <string>itpe-topic-splitter</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/turtlesoup0-macmini/Library/Logs/itpe-splitter/cloudflared.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/turtlesoup0-macmini/Library/Logs/itpe-splitter/cloudflared.err.log</string>
</dict>
</plist>
