# deploy/launchd/com.itpe.splitter.uvicorn.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.itpe.splitter.uvicorn</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>web.app:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8080</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/turtlesoup0-macmini/Projects/itpe-topic-splitter</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>LLM_PROVIDER</key>
        <string>mlx</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/Users/turtlesoup0-macmini/Library/Logs/itpe-splitter/uvicorn.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/turtlesoup0-macmini/Library/Logs/itpe-splitter/uvicorn.err.log</string>
</dict>
</plist>
